# epaper-frame.service
[Unit]
Description=E-Paper Photo Frame Web Server
After=network.target

[Service]
Type=simple
User=r2
WorkingDirectory=/home/r2/epaper-frame
ExecStart=/usr/bin/python3 /home/r2/epaper-frame/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
