# epaper-rotate.service
[Unit]
Description=FrameSync Image Rotation
After=network.target

[Service]
Type=oneshot
User=r2
WorkingDirectory=/home/r2/frame-sync
ExecStart=/usr/bin/python3 /home/r2/frame-sync/rotate_image.py
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
