# epaper-rotate.service
[Unit]
Description=E-Paper Photo Frame Image Rotation
After=network.target

[Service]
Type=oneshot
User=r2
WorkingDirectory=/home/r2/epaper-frame
ExecStart=/usr/bin/python3 /home/r2/epaper-frame/rotate_image.py
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
