# epaper-frame.service
[Unit]
Description=FrameSync Web Server
After=network.target

[Service]
Type=simple
User=r2
WorkingDirectory=/home/r2/frame-sync
ExecStart=/usr/bin/python3 /home/r2/frame-sync/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
